# crontab
*/5	*	*	*	*	~/Documents/eqmac-notify/eqmac-notify.sh
